# rc.local
#!/bin/sh
# /data/rc.local – simple one-shot launcher for install-janitza.sh

# Launch the installer in the background, redirecting its output to the log
nohup sh /data/install-janitza.sh >> /data/janitza/install-janitza.log 2>&1 &

# Exit right away so the boot process is not held up
exit 0
